# postinstall_servers_configurations/roles/akfgen/files/akfgen.timer
[Unit]
Description=AKFGEN Timer

[Timer]
OnBootSec=0
OnUnitActiveSec=1h
Unit=akfgen.service

[Install]
WantedBy=multi-user.target